# root/etc/services.d/python/run
#!/usr/bin/with-contenv sh
WMT_CONFIG=hiking python3 /waymarkedtrails/frontend.py
