# root/etc/services.d/python/run
#!/usr/bin/with-contenv sh
python3 /waymarkedtrails/frontend.py
